# server-files/supervisor.conf
[program:sscaas]
command = /var/www/go/bin/sscaas
user = sscaas
stdout_logfile = /var/log/gunicorn_supervisor.log
redirect_stderr = true
